# slett-opptak.service
[Unit]
Description=Slett gamle videoopptak

[Service]
Type=simple
User=opptaker
ExecStart=/usr/local/bin/slettopptak
